# etc/rc.local
#!/usr/local/bin/bash

(/lib/smartdc/joyent_rc.local)

exit 0
